# utilities/sync-screenshots/systemd/sync-screenshots-linker.timer
[Unit]
Description=Timer for Steam screenshot sync symlinker service

[Timer]
OnBootSec=60s
OnUnitActiveSec=60s
Unit=sync-screenshots-linker.service
