# utilities/steam-screenshot-sync/systemd/sync-screenshots.path
[Unit]
Description=Sync Steam Screenshots

[Path]
PathModified=%h/.steam_screenshots
Unit=sync-screenshots.service

[Install]
WantedBy=default.target
