# utilities/sync-screenshots/sync-screenshots.service
[Unit]
Description=Sync Steam Screenshots

[Service]
Type=oneshot
ExecStart=%h/SteamOS-Tools/utilties/sync-screenshots/sync-screenshots.sh
