# cfgs/systemd/ludusavi-backup.timer
[Unit]
Description="Ludusavi backup timer"

[Timer]
OnCalendar=*-*-* *:00/5:00
Unit=ludusavi-backup.service


[Install]
WantedBy=timers.target
