# cfgs/systemd/ludusavi-backup.service
[Unit]
Description="Ludusavi backup"

[Service]
ExecStart=/opt/ludusavi backup --force
