# cfgs/systemd/ludusavi-backup.timer
[Unit]
Description="Ludusavi backup timer"

[Timer]
OnCalendar=*/30 * * * *
Unit=ludusavi-backup.service

[Install]
WantedBy=timers.target
